# debian/etc/cron.daily/primming-certbot
#!/bin/sh -e

/usr/bin/certbot renew && cd /opt/primming-docker/ && /usr/bin/docker-compose -p primming restart proxy